# server/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
directory=/usr/src/server/
command=python3 server.py
autostart=true
autorestart=true
priority=100
stdout_logfile=/usr/src/server/python_test_output.log